# encrypt.service
[Unit]
Description=Encrypt Decrypt Service

[Service]
ExecStart=/home/uzmah/src/Service_TP/check_encrypt.sh
Restart=always
